# supervisor/etc/supervisor/conf.d/graftnoded.conf
[program:graftnoded]
#environment=DISPLAY=:20
command=/opt/graftnoded --testnet --detach
priority=20
user=root
#directory=/home/automation
#environment=HOME="/home/automation"
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stdout_logfile_backups=4
stdout_logfile_maxbytes=10MB
stderr_logfile=NONE
autostart=true
autorestart=true
startretries=0
numprocs=1